# mine.target
[Unit]
Description=mine: User target
After=xorg.target

[Install]
Alias=default.target
