# wm.target
[Unit]
Description=wm: Window manager
Wants=xorg.target
Wants=mine.target
Requires=dbus.socket
AllowIsolate=true

[Install]
Alias=default.target
