# xsession.target
[Unit]
Description=Manually started target to denote that an X session has started
BindsTo=graphical-session.target
